# systemd/user/claude-bun-zombie-cleaner.service
[Unit]
Description=Clean up orphan plugin:telegram bun processes with deleted cwd

[Service]
Type=oneshot
ExecStart=/home/gmoh/dotfiles/bin/claude-bun-zombie-cleaner
Nice=10
