# systemd/user/claude-control-bot.service
[Unit]
Description=Claude Control Bot (Telegram → cs commands)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment="PATH=/home/gmoh/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/gmoh/.local/bin/claude-control-bot
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/gmoh/.claude/logs/claude-control-bot.log
StandardError=append:/home/gmoh/.claude/logs/claude-control-bot.log

[Install]
WantedBy=default.target
