# systemd/user/claude-memory-lifecycle.timer
[Unit]
Description=Claude Memory Lifecycle Timer (매일 03:30 KST)

[Timer]
OnCalendar=*-*-* 03:30:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
