# systemd/user/claude-bun-ensure.service
[Unit]
Description=Claude Code 부팅 후 bun (telegram MCP) 누락 세션 자동 복구
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/gmoh/.local/bin/claude-bun-ensure
StandardOutput=append:/home/gmoh/.claude/logs/bun-ensure.log
StandardError=append:/home/gmoh/.claude/logs/bun-ensure.log

[Install]
WantedBy=default.target
