# systemd/user/claude-new-loc.service
[Unit]
Description=Claude Code session for new_loc (telegram remote-control)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/gmoh/av-ros/test/new_loc
Environment="PATH=/home/gmoh/.bun/bin:/home/gmoh/.nvm/versions/node/v22.17.0/bin:/home/gmoh/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="TELEGRAM_STATE_DIR=/home/gmoh/av-ros/test/new_loc/.claude/telegram"
Environment="TERM=xterm-256color"
ExecStart=/home/gmoh/.local/bin/claude-service-exec
ExecStartPost=/home/gmoh/.local/bin/claude-service-notify start claude-new-loc
ExecStopPost=/home/gmoh/.local/bin/claude-service-notify stop claude-new-loc
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/gmoh/.claude/logs/claude-new-loc.log
StandardError=append:/home/gmoh/.claude/logs/claude-new-loc.log
KillMode=mixed
TimeoutStopSec=20

[Install]
WantedBy=default.target
