# systemd/user/claude-resume.service
[Unit]
Description=Claude Code session for resume (telegram remote-control)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/gmoh/resume
Environment="PATH=/home/gmoh/.bun/bin:/home/gmoh/.nvm/versions/node/v22.17.0/bin:/home/gmoh/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="TELEGRAM_STATE_DIR=/home/gmoh/resume/.claude/telegram"
Environment="TERM=xterm-256color"
ExecStart=/home/gmoh/.local/bin/claude-service-exec
ExecStartPost=/home/gmoh/.local/bin/claude-service-notify start claude-resume
ExecStopPost=/home/gmoh/.local/bin/claude-service-notify stop claude-resume
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/gmoh/.claude/logs/claude-resume.log
StandardError=append:/home/gmoh/.claude/logs/claude-resume.log
KillMode=mixed
TimeoutStopSec=20

[Install]
WantedBy=default.target
