# systemd/user/claude-raion.service
[Unit]
Description=Claude Code session for raion (telegram remote-control)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/home/gmoh/raion/.claude/telegram/.env

[Service]
Type=simple
WorkingDirectory=/home/gmoh/raion
Environment="CLAUDE_SESSION_NAME=raion"
Environment="PATH=/home/gmoh/.bun/bin:/home/gmoh/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="TELEGRAM_STATE_DIR=/home/gmoh/raion/.claude/telegram"
Environment="TERM=xterm-256color"
ExecStart=/home/gmoh/.local/bin/claude-service-exec
ExecStartPost=/home/gmoh/.local/bin/claude-service-notify start claude-raion
ExecStopPost=/home/gmoh/.local/bin/claude-service-notify stop claude-raion
Restart=always
RestartSec=5
StandardOutput=append:/home/gmoh/.claude/logs/claude-raion.log
StandardError=append:/home/gmoh/.claude/logs/claude-raion.log
KillMode=mixed
TimeoutStopSec=20

[Install]
WantedBy=default.target
