# systemd/user/claude-api-5xx-watcher.service
[Unit]
Description=Watch claude session jsonls for API 5xx errors and alert via telegram

[Service]
Type=oneshot
ExecStart=/home/gmoh/dotfiles/bin/claude-api-5xx-watcher
Nice=10
