# systemd/user/raion-auth-watch.service
[Unit]
Description=Raion claude.ai self-auth health probe (로그인 만료 감지 → todobot 알림)
After=network-online.target
# raion 전용 — todobot 토큰이 있는 머신에서만 기동
ConditionPathExists=%h/raion/todo-sync/bot.token

[Service]
Type=oneshot
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=%h/.local/bin/raion-claude-auth-watch
StandardOutput=append:%h/.claude/logs/raion-auth-watch.log
StandardError=append:%h/.claude/logs/raion-auth-watch.log
